# deploy/systemd/phd-capital-backend.service
[Unit]
Description=PHD Capital Rationale Studio Backend
After=network.target postgresql.service

[Service]
Type=simple
User=phd
Group=phd
WorkingDirectory=/home/phd/app
EnvironmentFile=/etc/phd-capital.env

ExecStart=/home/phd/app/venv/bin/gunicorn \
    --bind 127.0.0.1:5000 \
    --workers 4 \
    --timeout 300 \
    --access-logfile /var/log/phd-capital/access.log \
    --error-logfile /var/log/phd-capital/error.log \
    --log-level info \
    backend.app:app

Restart=always
RestartSec=10

StandardOutput=append:/var/log/phd-capital/backend.log
StandardError=append:/var/log/phd-capital/backend-error.log

[Install]
WantedBy=multi-user.target
